# docs/source/advice/15/celeryd.conf
[program:courses_celeryd]
command=/home/zdimon/co_ve/co/manage.py celeryd
directory=/home/zdimon/co_ve/co
environment=PATH="/home/zdimon/co_ve/bin"
user=zdimon
autostart=true
autorestart=true
